# scripts/vpn-nat/dns-update.service
[Unit]
Description=Update VPN/NAT DNS Record on Boot
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
ExecStart=/usr/local/bin/update-dns-on-boot.sh
TimeoutStartSec=300
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
